# funke/systemd/funke-indexer.service
[Unit]
Description=Funke app and file indexer

[Service]
Type=oneshot
ExecStart=/

